# roll_one.cron
*/15 * * * * /home/pi/robots/roll_one_for_me/revive.sh
